# usr/lib/systemd/system/tftp.socket
[Unit]
Description=Tftp Server Activation Socket

[Socket]
ListenDatagram=172.16.41.244:69

[Install]
WantedBy=sockets.target
