# usr/lib/systemd/system/tftp.socket
[Unit]
Description=Tftp Server Activation Socket

[Socket]
ListenDatagram=BASTION_IP:69

[Install]
WantedBy=sockets.target
